# resources/io.github.trbjo.bob.launcher.service
[Unit]
Description=BobLauncher
ConditionEnvironment=WAYLAND_DISPLAY
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Slice=app.slice
Type=exec
#EnvironmentFile=%h/.config/environment
ExecStart=%h/.local/bin/bob-launcher --hidden

[Install]
WantedBy=graphical-session.target
